# resources/etc/supervisor/conf.d/mariadb.conf
[program:mariadb]
directory=/var/lib/mysql/
user=mysql
command=docker-entrypoint.sh mariadbd
autostart=true
autorestart=false
redirect_stderr=true
startsecs=0
